# fjosetinfo.service
# FjosetInfo systemd service file for Raspberry Pi
# 
# Usage:
# 1. Copy this file to /etc/systemd/system/fjosetinfo.service
# 2. Edit the paths and user as needed
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable fjosetinfo
# 5. Run: sudo systemctl start fjosetinfo

[Unit]
Description=FjosetInfo Display Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/FjosetInfo
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fjosetinfo

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=8787

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/pi/FjosetInfo

[Install]
WantedBy=multi-user.target